# ops/systemd/minecraft-local-backup.service
[Unit]
Description=Minecraft nightly local archive backup
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
# Update this path if the repo lives somewhere other than /opt/minecraft-backup-google-drive
WorkingDirectory=/opt/minecraft-backup-google-drive
ExecStart=/opt/minecraft-backup-google-drive/scripts/nightly-local-backup.sh

