# ops/systemd/minecraft-local-backup.timer
[Unit]
Description=Run Minecraft nightly local archive backup

[Timer]
OnCalendar=*-*-* 03:30:00
Persistent=true

[Install]
WantedBy=timers.target

